# systemd/amnezia-monitor-collector.service
[Unit]
Description=AmneziaWG traffic monitor collector
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/amnezia-monitor
ExecStart=/opt/amnezia-monitor/.venv/bin/python -m app.collector --config /opt/amnezia-monitor/config.toml
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
